# foreground-booster.service.in
[Unit]
Description=Foreground App CPU weight booster

[Service]
ExecStart=@KDE_INSTALL_FULL_BINDIR@/bin/foreground_booster
