# foreground-booster.service.in
[Unit]
Description=Foreground app CPU weight booster
PartOf=graphical-session.target
After=graphical-session.target

[Service]
Type=simple
Restart=on-failure
ExecStart=@KDE_INSTALL_FULL_BINDIR@/foreground_booster
Slice=background.slice

[Install]
WantedBy=graphical-session.target
