# plasma-foreground-booster.service.in
[Unit]
Description=Foreground app CPU weight booster
PartOf=graphical-session.target
After=graphical-session.target

[Service]
Type=dbus
BusName=org.kde.foreground-booster
Restart=on-failure
ExecStart=@KDE_INSTALL_FULL_BINDIR@/foreground_booster
ExecCondition=/usr/bin/kde-systemd-start-condition --condition "kcgroupsrc:Foreground Booster:autostart:false"
Slice=background.slice

[Install]
WantedBy=graphical-session.target
